# systemd/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway (v2026.3.8)
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/node /home/noestelar/.npm-global/lib/node_modules/openclaw/dist/entry.js gateway --port 18789
Restart=always
RestartSec=5
TimeoutStopSec=30
TimeoutStartSec=30
SuccessExitStatus=0 143
KillMode=control-group
Environment=HOME=/home/noestelar
EnvironmentFile=/home/noestelar/.openclaw/credentials/1password.env
Environment=TMPDIR=/tmp
Environment=PATH=/home/noestelar/.local/bin:/home/noestelar/.npm-global/bin:/home/noestelar/bin:/home/noestelar/.volta/bin:/home/noestelar/.asdf/shims:/home/noestelar/.bun/bin:/home/noestelar/.nvm/current/bin:/home/noestelar/.fnm/current/bin:/home/noestelar/.local/share/pnpm:/usr/local/bin:/usr/bin:/bin
Environment=OPENCLAW_GATEWAY_PORT=18789
Environment=OPENCLAW_SYSTEMD_UNIT=openclaw-gateway.service
Environment="OPENCLAW_WINDOWS_TASK_NAME=OpenClaw Gateway"
Environment=OPENCLAW_SERVICE_MARKER=openclaw
Environment=OPENCLAW_SERVICE_KIND=gateway
Environment=OPENCLAW_SERVICE_VERSION=2026.3.8

[Install]
WantedBy=default.target
